# deploy/ai-challenge-api.service
[Unit]
Description=Thử thách AI submissions API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ai-challenge-hub
EnvironmentFile=/etc/ai-challenge-hub.env
Environment=NODE_NO_WARNINGS=1
ExecStart=/usr/bin/node /opt/ai-challenge-hub/server/index.mjs
Restart=always
RestartSec=5
User=aichallenge
Group=aichallenge
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0077
ReadWritePaths=/opt/ai-challenge-hub/storage

[Install]
WantedBy=multi-user.target
