# deploy/ai-challenge-backup.service
[Unit]
Description=Thử thách AI production backup
Wants=ai-challenge-api.service

[Service]
Type=oneshot
WorkingDirectory=/opt/ai-challenge-hub
EnvironmentFile=/etc/ai-challenge-hub.env
Environment=NODE_NO_WARNINGS=1
ExecStart=/usr/bin/node /opt/ai-challenge-hub/scripts/backup-production.mjs
User=aichallenge
Group=aichallenge
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0077
ReadWritePaths=/opt/ai-challenge-hub/storage
